# supervisord.conf
[supervisord]
nodaemon=true

[program:startxfce4]
priority=10
directory=/root
command=/usr/bin/startxfce4
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/xfce4.log
stderr_logfile=/var/log/xfce4.err

[program:xvfb]
priority=10
directory=/
command=/usr/bin/Xvfb :1 -screen 0 1024x768x16
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:x11vnc]
priority=10
directory=/
command=x11vnc -display :1 -xkb -forever -loop -viewonly -noxdamage -repeat -rfbauth /root/.vnc/user_password -rfbport 5900 -shared
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/x11vnc-user.log
stderr_logfile=/var/log/x11vnc-user.err

[program:x11vnc-admin]
priority=10
directory=/
command=x11vnc -display :1 -xkb -forever -loop -noxdamage -repeat -rfbauth /root/.vnc/admin_password -rfbport 5901 -shared
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/x11vnc-admin.log
stderr_logfile=/var/log/x11vnc-admin.err

[program:novnc]
priority=20
environment=VNC_PORT=5900,NO_VNC_PORT=6900
command=/root/novnc_startup.sh
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/novnc-user.log
stderr_logfile=/var/log/novnc-user.err

[program:novnc-admin]
priority=20
environment=VNC_PORT=5901,NO_VNC_PORT=6901
command=/root/novnc_startup.sh
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/novnc-admin.log
stderr_logfile=/var/log/novnc-admin.err

[program:game]
priority=5
command=/root/game_startup.sh
autostart=true
stopsignal=QUIT
stdout_logfile=/var/log/game.log
stderr_logfile=/var/log/game.err
